# coffeebreak/templates/systemd.service.j2
# CoffeeBreak {{ service_name }} Service
# Generated by CoffeeBreak CLI for {{ domain }}
# Generated on: {{ timestamp }}

[Unit]
Description=CoffeeBreak {{ service_description | default(service_name) }}
Documentation=https://docs.coffeebreak.com/
After=network.target
{% if service_dependencies %}
{% for dep in service_dependencies %}
After={{ dep }}.service
Requires={{ dep }}.service
{% endfor %}
{% endif %}
{% if service_name == 'api' %}
After=postgresql.service
Requires=postgresql.service
{% endif %}

[Service]
Type={{ service_type | default('simple') }}
User={{ service_user | default('coffeebreak') }}
Group={{ service_group | default('coffeebreak') }}
WorkingDirectory={{ service_workdir | default('/opt/coffeebreak/' + service_name) }}
ExecStart={{ service_exec_start }}
{% if service_exec_reload %}
ExecReload={{ service_exec_reload }}
{% endif %}
{% if service_exec_stop %}
ExecStop={{ service_exec_stop }}
{% else %}
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec={{ service_stop_timeout | default('30') }}
{% endif %}

# Process management
Restart={{ service_restart | default('always') }}
RestartSec={{ service_restart_sec | default('5') }}
StartLimitInterval={{ service_start_limit_interval | default('60') }}
StartLimitBurst={{ service_start_limit_burst | default('5') }}

# Resource limits
{% if service_memory_limit %}
MemoryLimit={{ service_memory_limit }}
{% endif %}
{% if service_cpu_quota %}
CPUQuota={{ service_cpu_quota }}
{% endif %}

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ service_rw_paths | default(['/opt/coffeebreak/' + service_name + '/data', '/var/log/coffeebreak']) | join(' ') }}
{% if service_name != 'nginx' %}
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
{% endif %}

# Network security
{% if service_restrict_network | default(false) %}
PrivateNetwork=true
{% else %}
IPAddressDeny=any
IPAddressAllow=localhost
{% if service_allowed_ips %}
{% for ip in service_allowed_ips %}
IPAddressAllow={{ ip }}
{% endfor %}
{% endif %}
{% endif %}

# File system security
PrivateDevices=true
{% if service_name != 'nginx' %}
ProtectKernelLogs=true
ProtectClock=true
{% endif %}

# Capabilities
{% if service_capabilities %}
AmbientCapabilities={{ service_capabilities | join(' ') }}
CapabilityBoundingSet={{ service_capabilities | join(' ') }}
{% else %}
CapabilityBoundingSet=~CAP_SYS_ADMIN
{% endif %}

# Environment
Environment=NODE_ENV=production
Environment=LOG_LEVEL={{ service_log_level | default('info') }}
{% if service_environment %}
{% for env_var, env_value in service_environment.items() %}
Environment={{ env_var }}={{ env_value }}
{% endfor %}
{% endif %}

# Environment files
{% if service_env_files %}
{% for env_file in service_env_files %}
EnvironmentFile={{ env_file }}
{% endfor %}
{% else %}
EnvironmentFile=-/etc/coffeebreak/{{ service_name }}.env
EnvironmentFile=-/opt/coffeebreak/{{ service_name }}/.env
{% endif %}

# Logging
StandardOutput={{ service_stdout | default('journal') }}
StandardError={{ service_stderr | default('journal') }}
SyslogIdentifier=coffeebreak-{{ service_name }}

{% if service_name == 'nginx' %}
# Nginx specific settings
ExecStartPre=/usr/sbin/nginx -t
ExecReload=/bin/kill -s HUP $MAINPID
PIDFile=/var/run/nginx.pid
{% endif %}

{% if service_name == 'api' %}
# API specific settings
TimeoutStartSec=60
LimitNOFILE=65536
{% endif %}

[Install]
WantedBy=multi-user.target
{% if service_also_enable %}
Also={{ service_also_enable | join(' ') }}
{% endif %}